# overlay/etc/profile.d/common-env.sh
#!/bin/bash

export USER=${USER:-$(id -un)}
export HOME=${HOME:-~}
export TERM=${TERM:-linux}
